# systemd/backup-snapshots.service
[Unit]
Description=Local backup snapshots

[Service]
Type=oneshot
ExecStart=/home/gideon/.local/bin/take-snapshot
