# systemd/backup-snapshots.service
[Unit]
Description=Local backup snapshots

[Service]
Type=oneshot
ExecStart=/usr/local/bin/take-snapshot
