# systemd/dropbox.service
[Unit]
Description=Dropbox file synchronisation app

[Service]
ExecStart=/usr/bin/dropbox
