# systemd/duplicity.timer
[Unit]
Description=Duplicity Backup Timer

[Timer]
OnBootSec=10m
OnUnitActiveSec=120m
Unit=duplicity.service

[Install]
WantedBy=timers.target
