# systemd/duplicity.service
[Unit]
Description=Duplicity Backup Service

[Service]
Type=oneshot
ExecStart=/home/gideon/.local/bin/run-backup
