# systemd/freshclam.service
[Unit]
Description=Virus database update

[Service]
Type=oneshot
ExecStart=/usr/bin/freshclam
