# systemd/geoclue2.service
[Unit]
Description=Geoclue2 Location Updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/geoclue2-location-updater.sh
