# systemd/iwgtk-indicator.service
[Unit]
Description=IWD Tray Indicator

[Service]
ExecStart=iwgtk -i
