# systemd/keepassxc.service
[Unit]
Description=KeepassXC keyring manager

[Service]
ExecStart=/usr/bin/keepassxc
