# systemd/mbsync.service
[Unit]
Description=Mailbox Syncronisation Service

[Service]
Type=oneshot
ExecStart=/home/gideon/.local/bin/email-sync sync_remote
