# systemd/nextcloud-sync.service
[Unit]
Description=Nextcloud Synchronisation Service

[Service]
Type=oneshot
ExecStart=/home/gideon/.local/bin/nextcloud-sync
