# systemd/playerctld.service
[Unit]
Description=Media Player Control Daemon

[Service]
ExecStart=/usr/bin/playerctld
Restart=always

[Install]
WantedBy=default.target
