# systemd/vdirsyncer.service
[Unit]
Description=DAV Synchronisation Service (Calendar and Contacts)

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer -v ERROR sync
