# systemd/vdirsyncer.service
[Unit]
Description=Calendar Synchronisation Service

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer -v ERROR sync
